# systemd/llama-fast.service
[Unit]
Description=llama.cpp Vulkan - Qwen3-8B (port 8002)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brendanhouck
Restart=always
RestartSec=30
TimeoutStopSec=60
ExecStart=/home/brendanhouck/start_fast.sh
StandardOutput=append:/tmp/llama-fast.log
StandardError=append:/tmp/llama-fast.log

[Install]
WantedBy=multi-user.target
