# friday_rollout/systemd/vllm-coder.service
[Unit]
Description=vLLM Qwen3-Coder-30B-A3B (Tier 2 code) on card 0 TP=1
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=600
TimeoutStopSec=30
ExecStartPre=-/usr/bin/docker stop vllm-b70-coder
ExecStartPre=-/usr/bin/docker rm vllm-b70-coder
ExecStart=/usr/bin/docker run -d --privileged --shm-size 16g --net=host --ipc=host --restart=unless-stopped --device /dev/dri:/dev/dri -v /mnt/models:/llm/models --name=vllm-b70-coder --entrypoint= vllm-xpu:gemma4-fixed sleep infinity
ExecStartPost=/home/brendanhouck/start_coder.sh
ExecStop=/usr/bin/docker stop -t 15 vllm-b70-coder

[Install]
WantedBy=multi-user.target
